# system_files/usr/lib/systemd/system/armada-controller-type.service
[Unit]
Description=Apply Armada controller type
ConditionPathExists=/etc/armada/controller.conf
After=inputplumber.service
Requires=inputplumber.service
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/armada/controller-type apply

[Install]
WantedBy=multi-user.target
