# system_files/usr/lib/systemd/system/armada-powerd.service
[Unit]
Description=armada power policy daemon
StartLimitIntervalSec=0

[Service]
Type=dbus
BusName=org.armada.Power
NotifyAccess=main
ExecStart=/usr/libexec/armada/armada-powerd
ExecStopPost=/usr/libexec/armada/armada-powerd --fan-safe
Restart=always
RestartSec=1
WatchdogSec=15

[Install]
WantedBy=multi-user.target
